# config/supervisor-mysql.ini
[program:mysql]
command=mysqld_safe --user=%(ENV_MYSQL_USER)s --datadir=%(ENV_MYSQL_DATA_DIR)s
priority=1
autostart=true
autorestart=true
stdout_logfile=%(ENV_MYSQL_STDOUT_FILE)s
stderr_logfile=%(ENV_MYSQL_STDERR_FILE)s
